# server/config/supervisord.conf
[unix_http_server]
file=/tmp/supervisor.sock
chmod=0700

[inet_http_server]
port=0.0.0.0:9000

[supervisord]
logfile=./logs/supervisord.log
logfile_maxbytes=50MB
logfile_backups=10
loglevel=info
pidfile=/tmp/supervisord.pid
nodaemon=false
silent=false
minfds=1024
minprocs=200

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[program:service_price_predict]
command=gunicorn -c ./config/gunicorn.py app.app:app --preload
autostart=true
startsecs=1
autorestart=true
stopwaitsecs=1

[include]
files=./config/supervisord.conf